# config/crontab
# =============================================
#  misskey backup 
#  バックアップを自動実行する時間を定義します。
# =============================================
0 */12 * * * . /root/backup.sh